# rpi/call_status_heartbeat_consumer.service
[Unit]
Description=Call Status ESP32 heartbeat polling service

[Service]
Type=simple
ExecStart=/usr/bin/call_status.sh poll_heartbeats
Restart=always
RestartSec=2s

[Install]
WantedBy=default.target
